# deploy/render/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/app/logs/supervisord.log
pidfile=/var/run/supervisord.pid

[program:uvicorn]
command=uvicorn src.api.main:app --host 127.0.0.1 --port 8000 --workers 1 --loop uvloop
directory=/app
autostart=true
autorestart=true
startretries=3
stderr_logfile=/app/logs/api.err.log
stdout_logfile=/app/logs/api.out.log
environment=PATH="/root/.local/bin:%(ENV_PATH)s",PYTHONPATH="/app",ENVIRONMENT="production"

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
startretries=3
stderr_logfile=/app/logs/nginx.err.log
stdout_logfile=/app/logs/nginx.out.log